# service/okawak_blog-aws-credentials.timer
[Unit]
Description=Daily refresh of okawak Blog AWS runtime credentials

[Timer]
OnBootSec=10min
OnCalendar=*-*-* 04:05:00
Persistent=true
Unit=okawak_blog-aws-credentials.service

[Install]
WantedBy=timers.target
